# rootfs/etc/s6-overlay/s6-rc.d/cups-server/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

bashio::log.info "Starting CUPS server from S6"

retries=0
until [ -e /var/run/avahi-daemon/socket ]; do
  if [ "$retries" -ge 30 ]; then
    bashio::log.error "Timed out waiting for Avahi socket (/var/run/avahi-daemon/socket); aborting."
    exit 1
  fi
  retries=$((retries + 1))
  sleep 1s
done

cupsd -f

