# rootfs/etc/s6-overlay/s6-rc.d/avahi-daemon/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

bashio::log.info "Starting Avahi daemon from S6"

WAIT_MAX_RETRIES=30
WAIT_INTERVAL=1
retries=0
until [ -e /run/dbus/system_bus_socket ]; do
  if [ "$retries" -ge "$WAIT_MAX_RETRIES" ]; then
    bashio::log.error "Timed out waiting for D-Bus socket (/run/dbus/system_bus_socket); aborting."
    exit 1
  fi
  retries=$((retries + 1))
  sleep "${WAIT_INTERVAL}s"
done

exec avahi-daemon
