# rootfs/etc/s6-overlay/s6-rc.d/initialization/up
#!/usr/bin/with-contenv bashio
set -euo pipefail

mkdir -p /var/run/dbus

