# rootfs/etc/s6-overlay/s6-rc.d/dbus-daemon/run
#!/usr/bin/with-contenv bashio

ulimit -n 1048576 2>/dev/null || true

bashio::log.info "Starting DBUS daemon from S6"

dbus-daemon --system --nofork
