# rootfs/etc/s6-overlay/s6-rc.d/dbus-daemon/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

bashio::log.info "Starting DBUS daemon from S6"

dbus-daemon --system --nofork
