# rootfs/etc/s6-overlay/s6-rc.d/cups-browsed/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting cups-browsed (AirPrint advertisement) from S6"

exec cups-browsed --no-daemon
