# rootfs/etc/s6-overlay/s6-rc.d/ingress-server/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

bashio::log.info "Starting ingress web server..."
exec nginx -g 'daemon off;'
