# rootfs/etc/s6-overlay/s6-rc.d/cups-config/up
/usr/bin/with-contenv bashio /usr/share/cups-config.sh
