# rootfs/etc/s6-overlay/s6-rc.d/cups-config/up
/usr/bin/with-contenv /bin/bash /usr/share/cups-config.sh
